# server/verda-sxtelo.socket
[Unit]
Description=Activation socket for Verda Ŝtelo

[Socket]
ListenStream=5142

[Install]
WantedBy=sockets.target
